# supervisord-init.conf
[program:init]
command:sh /bootstrap.sh
startsecs = 0
autorestart = false
startretries = 1